# grive-sync@.service
[Unit]
Description=Grive daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/grive-sync

[Install]
WantedBy=default.target
